# components/scale-radio-starter/payload/current_dev/systemd-dropins/plymouth-quit.service.d/mediastreamer-direct-now-playing-topmost-keeper.conf
[Service]
ExecStartPre=/bin/bash -lc 'for i in $(seq 1 240); do [ -f /tmp/mediastreamer-kiosk-ready ] && exit 0; sleep 0.2; done; exit 0'
